# com.ikuwow.minecraft-playtime-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ikuwow.minecraft-playtime-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>exec "$HOME/.local/bin/minecraft_playtime_monitor.py" --weekly-limits 90,90,90,90,90,120,120 --curfew-hour 22 --notify-remaining 60,10,5,1 --notify-repeat-interval 10 --process-pattern "java.*minecraft"</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
